# deployment/algoritmika-web.service
[Unit]
Description=Algoritmika Web Application
After=network.target postgresql.service
Wants=postgresql.service

[Service]
User=algoritmika
Group=algoritmika
WorkingDirectory=/opt/algoritmika
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="FLASK_APP=main.py"
Environment="FLASK_ENV=production"
Environment="PYTHONPATH=/opt/algoritmika"
EnvironmentFile=/opt/algoritmika/.env
ExecStart=/usr/bin/python3 -m gunicorn --config /opt/algoritmika/deployment/gunicorn.conf.py main:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target